# msc/reportcard-supervisor.conf
[program:reportcard]
command=/srv/reportcard/lib/server/server.co
directory=/srv/reportcard
user=www-data
environment=NODE_ENV=prod
